# .config/systemd/user/notes_git_sync.service
[Unit]
Description=Synchronizes notes

[Service]
Type=simple
ExecStart=fish -c "notes_git_sync"
Restart=always

[Install]
WantedBy=default.target
