# systemd/user/navi-update.service
[Unit]
Description=Update Navi repos

[Service]
Type=oneshot
ExecStart=/home/tim/.config/systemd/user/navi-update.fish
